# armv8/base-files/etc/init.d/drm-panfrost
#!/bin/sh /etc/rc.common

chmod -R 777 /dev/dri
